# task2/crontab
30 18 * * 5 bash /home/aleksandra/Desktop/Task8Bash/Task2 
